# mongoose.service.conf
[program:mongo]
command=/usr/bin/python /sleepy.mongoose/httpd.py --xorigin 
